# arkbox/arkbox/rpi/udev/99-wifi.rules
KERNEL=="wlan*", RUN+="/bin/sh -c /opt/arkbox/rpi/bin/starter.sh"
