# arkbox/arkbox/rpi/services/arkbox.service
[Unit]
Description=ArkBox Service

[Service]
Type=idle
RemainAfterExit=yes
ExecStart=/opt/arkbox/init.d/arkbox_alt start
ExecStop=/opt/arkbox/init.d/arkbox_alt stop
#Restart=/opt/arkbox/init.d/arkbox_alt restart
#PIDFile=/var/run/arkbox.pid

[Install]
WantedBy=multi-user.target
